# deploy/bdoom-api.service
[Unit]
Description=BDoom NestJS API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/bdoom/app
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3000
Environment=BDOOM_DB_PATH=/opt/bdoom/data/bdoom.sqlite
EnvironmentFile=/opt/bdoom/.env
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
